# autostart/.config/autostart/plank.desktop
[Desktop Entry]
Type=Application
Exec=plank
Hidden=false
X-MATE-Autostart-enabled=true
Name[en_GB]=Plank
Name=Plank
Comment[en_GB]=
Comment=
